# src/mosquitto-1.4.5/service/svscan/run
#!/bin/sh

/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
