# src/mosquitto-1.4.5/service/upstart/mosquitto.conf
description "Mosquitto MQTTv3.1 broker"
author "Roger Light <roger@atchoo.org>"

start on net-device-up

respawn

exec /usr/local/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
